# sway/dot-config/systemd/user/gnome-keyring-daemon.service.d/dbus.conf
[Service]
Type=dbus
BusName=org.freedesktop.secrets
