# sway/dot-config/systemd/user/gnome-keyring-daemon.service.d/unlock.conf
[Service]
ExecStart=
ExecStart=/usr/bin/gnome-keyring-daemon --foreground --control-directory=%t/keyring --components=secrets --unlock
StandardInput=data
StandardInputData=TG9naW4=
