# dotfiles/.config/systemd/user/dotfiles-pull.service
[Unit]
Description=Pull dotfiles repo automatically
Documentation=https://github.com/YHNdnzj/dotfiles

[Service]
WorkingDirectory=%h/dotfiles
ExecStartPre=/usr/bin/nm-online -q
ExecStart=/usr/bin/git pull

[Install]
WantedBy=default.target
