# dotfiles/.config/systemd/user/dotfiles-pull.service
[Unit]
Description=Pull dotfiles repo automatically
Documentation=https://github.com/YHNdnzj/dotfiles

[Service]
WorkingDirectory=%h/dotfiles
ExecStartPre=/usr/bin/wait-online
ExecStart=/usr/bin/git pull
Restart=on-failure
RestartSec=20

[Install]
WantedBy=default.target
