# ssh/.config/systemd/user/ssh-agent.service.d/askpass.conf
[Service]
PassEnvironment=WAYLAND_DISPLAY DISPLAY SSH_ASKPASS
