# paru/.config/systemd/user/paru-cache-cleanup.timer
[Unit]
Description=Cleanup outdated paru AUR package cache weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
