# paru/.config/systemd/user/repo-sync-to-kbfs.path
[Unit]
Description=Sync Local repo to KBFS after each database update

[Path]
PathChanged=%h/repo/Local/x86_64/yhndnzj.db.tar.gz

[Install]
WantedBy=default.target
